# php-worker/supervisord.d/union.cats831.conf
[program:cats831]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/clients/cats831/artisan queue:work redis --timeout=1800 --sleep=3 --tries=3 --daemon --queue=cats831_jobs
autostart=true
autorestart=true
numprocs=1
user=laradock
redirect_stderr=true
stdout_logfile=/var/www/clients/cats831/storage/logs/queue-jobs.log
